# AAA Ira/stepperslife-complete-20251118-020001/configs/classes-stepperslife.service
[Unit]
Description=Classes Stepperslife Website
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/websites/classes-stepperslife

# Build verification before starting
#ExecStartPre=/usr/local/bin/verify-nextjs-build.sh /root/websites/classes-stepperslife

# Node.js path from nvm
Environment="PATH=/root/.nvm/versions/node/v22.19.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="NODE_ENV=production"
Environment="PORT=3009"

# Load environment variables from .env file if it exists
EnvironmentFile=-/root/websites/classes-stepperslife/.env
EnvironmentFile=-/root/websites/classes-stepperslife/.env.local
EnvironmentFile=-/root/websites/classes-stepperslife/.env.production

# Start command
ExecStart=/root/.nvm/versions/node/v22.19.0/bin/npm start

# Restart configuration with limits
Restart=on-failure
RestartSec=30s
StartLimitBurst=5
StartLimitIntervalSec=600s

# Timeouts
TimeoutStartSec=120s
TimeoutStopSec=30s

# Logging to journal (not files)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=classes-stepperslife

# Process management
KillMode=mixed
KillSignal=SIGTERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=512
CPUQuota=40%
MemoryMax=1.5G
MemoryHigh=1G
TasksMax=512

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/root/websites/classes-stepperslife

[Install]
WantedBy=multi-user.target
